# scripts/scholarport-frontend.service
[Unit]
Description=Scholarport Frontend Application
Documentation=https://github.com/yourusername/scholarport-frontend
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/ubuntu/scholarport-frontend
ExecStart=/usr/local/bin/docker-compose -f docker-compose.prod.yml up -d
ExecStop=/usr/local/bin/docker-compose -f docker-compose.prod.yml down
ExecReload=/usr/local/bin/docker-compose -f docker-compose.prod.yml restart
User=ubuntu
Group=ubuntu
Restart=on-failure
RestartSec=10s

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=scholarport-frontend

[Install]
WantedBy=multi-user.target
